# card_reader_adaptor.conf
[program:device-controller]
user=root
command=/opt/local/card_reader_adaptor_bin
autostart = true
autorestart = true

